# crontab
*/15 * * * * * php -d display_errors=1 -d error_reporting=E_ALL /healthcheck/run.php
