# crontab
*/15 * * * * * php /healthcheck/run.php
